# Client/makerscreen.service
[Unit]
Description=MakerScreen Display Client
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/makerscreen
ExecStart=/usr/bin/python3 /opt/makerscreen/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
